# build/inittab
::sysinit:/bin/cat /etc/motd
::sysinit:/bin/hostname webOS

ttyS0::respawn:-/bin/sh

::ctrlaltdel:/sbin/reboot
::shutdown:/bin/umount -a -r
