# build/inittab
::sysinit:/bin/cat /etc/motd
::sysinit:/bin/hostname webOS
::sysinit:/bin/sh /etc/profile

ttyS0::respawn:/bin/sh -c "source /etc/profile && exec /bin/sh"

::ctrlaltdel:/sbin/reboot
::shutdown:/bin/umount -a -r
